# .config/systemd/user/nitrogen.timer
[Unit]
Description=Run nitrogen every 30 minutes

[Timer]
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
